# deploy/systemd/limeclicks-celerybeat.service
[Unit]
Description=LimeClicks Celery Beat Scheduler Service
After=network.target postgresql.service redis.service limeclicks-celery.service
Wants=postgresql.service redis.service limeclicks-celery.service

[Service]
Type=simple
User=limeclicks
Group=limeclicks
WorkingDirectory=/home/limeclicks/limeclicks
Environment="PATH=/home/limeclicks/.pyenv/shims:/home/limeclicks/.pyenv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYENV_ROOT=/home/limeclicks/.pyenv"

# Load environment variables from .env file
EnvironmentFile=/home/limeclicks/.env

# Celery Beat configuration
Environment="CELERY_APP=limeclicks"
Environment="CELERYBEAT_LOG_LEVEL=INFO"
Environment="CELERYBEAT_LOG_FILE=/var/log/limeclicks/celerybeat.log"
Environment="CELERYBEAT_SCHEDULE=/home/limeclicks/limeclicks/celerybeat-schedule"

# Start command
ExecStart=/home/limeclicks/limeclicks/venv/bin/celery -A ${CELERY_APP} beat \
          --loglevel=${CELERYBEAT_LOG_LEVEL} \
          --logfile=${CELERYBEAT_LOG_FILE} \
          --schedule=${CELERYBEAT_SCHEDULE}

# Restart policy
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/limeclicks/limeclicks /var/log/limeclicks

[Install]
WantedBy=multi-user.target